# ansible/roles/vault/templates/etc/profile.d/vault.sh.j2
#!/bin/bash
export VAULT_ADDR='https://{{ ansible_fqdn }}:{{ vault_port }}'
export VAULT_CACERT='{{ vault_certfile }}'
